# notebook/etc-config/socat-pveproxy.service
[Unit]
Description=socat relay: WiFi -> Proxmox LAN 192.168.1.40:8006
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/socat-pveproxy.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
